# vps-deploy/cenopie-backend.service
[Unit]
Description=Cenopie Backend Service
Documentation=https://github.com/your-username/cenopie
After=network.target mongod.service redis.service

[Service]
Type=simple
User=USER_PLACEHOLDER
WorkingDirectory=/opt/cenopie/backend
Environment=NODE_ENV=production
Environment=PORT=4000
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cenopie-backend

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/cenopie
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Health check
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target