# vps-deploy/cenopie-frontend.service
[Unit]
Description=Cenopie Frontend Service
Documentation=https://github.com/your-username/cenopie
After=network.target cenopie-backend.service

[Service]
Type=simple
User=USER_PLACEHOLDER
WorkingDirectory=/opt/cenopie/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cenopie-frontend

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/cenopie
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Health check
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target